# files/etc/crontabs/root
*/1 * * * * /opt/mesh/retrieve_ip
0 * * * * /opt/mesh/pw_reset
