# files/etc/crontabs/root
*/1 * * * * /opt/mesh/retrieve_ip
* * */1 * * rm /etc/shawdow 

