# extender-node/files/etc/uci-defaults/98_permissions
#!/bin/sh

# Fix dropbear permissions
chmod 0700 /etc/dropbear
chmod 0600 /etc/dropbear/authorized_keys

exit 0