# files/etc/uci-defaults/90_sudomesh_blank_homenode
#!/bin/sh

UUID=$(cat /proc/sys/kernel/random/uuid)
uci set tunneldigger.main.uuid=$UUID
uci commit tunneldigger

exit 0
