# files/etc/uci-defaults/90_sudomesh_blank_homenode
#!/bin/sh

uci set network.lan.ipaddr="172.22.0.1"
uci set network.lan.netmask="255.255.255.0"
uci set network.lan.proto="static"

uci commit network

exit 0
