# files/etc/uci-defaults/99_start_cron
/etc/init.d/cron start
/etc/init.d/cron enable
